# Source/Arch/Arch-Root/etc/systemd/system/pcie_aspm.timer
[Unit]
Description="Pcie_ASPM"

[Timer]
OnStartupSec=30

[Install]
WantedBy=multi-user.target
