# Source/Arch/Arch-Root/etc/systemd/system/pcie_aspm.service
[Unit]
Description=Pcie_ASPM

[Service]
ExecStart=/usr/bin/bash -c "echo powersupersave > /sys/module/pcie_aspm/parameters/policy"
Type=oneshot
